# ops/email/openclaw-gmail-pull.timer
[Unit]
Description=Run OpenClaw Gmail Pull every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=openclaw-gmail-pull.service
Persistent=true

[Install]
WantedBy=timers.target
